# dservice.cron
* * * * * /home/metastore/sotstore/dservice/watchdog.sh > ~/wd.log

